# ldap-samba/config/supervisor-updatesid.conf
[program:updatesdi]
#user=apache
command=sh /updateSID.sh
