# supervisor-updatesid.conf
[program:updatesdi]
command=sh /updateSID.sh
